# deploy/drone-delivery-server.service
[Unit]
Description=Drone Delivery Server (broker + server + sim drones)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/drone-delivery
ExecStart=/opt/drone-delivery/deploy/compose-wrapper.sh server up -d
ExecStop=/opt/drone-delivery/deploy/compose-wrapper.sh server down
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
