# deploy/drone-delivery-drone.service
[Unit]
Description=Drone Delivery Drone Client
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/drone-delivery
ExecStart=/opt/drone-delivery/deploy/compose-wrapper.sh drone up -d
ExecStop=/opt/drone-delivery/deploy/compose-wrapper.sh drone down
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
